# deploy/systemd/llm-wiki-slack-realtime.service
[Unit]
Description=LLM wiki Slack near-real-time incremental chunk sync
Documentation=file:%h/Work/Workspaces/displayr-llm-wiki/docs/slack-realtime-sync.md

[Service]
Type=oneshot
WorkingDirectory=%h/Work/Workspaces/displayr-llm-wiki
Environment=PATH=%h/.local/bin:%h/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=%h/Work/Workspaces/displayr-llm-wiki/scripts/run_slack_realtime_update.sh
